# news-cache.timer
[Unit]
Description=Update news cache every 15 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
